# scripts/systemd/snmpsimd.socket
[Unit]
Description=SNMP Simulator Socket Activation

[Socket]
ListenDatagram=1161

[Install]
WantedBy=multi-user.target
